# server/root/etc/services.d/app/run
#!/usr/bin/with-contenv sh

cd /opt/app
s6-setuidgid app python -m openapi_server